# config/com.windsurf.integration.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.windsurf.integration</string>
    <key>ProgramArguments</key>
    <array>
        <string>${PWD}/venv/bin/python3</string>
        <string>${PWD}/main.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>${PWD}/logs/windsurf_integration_error.log</string>
    <key>StandardOutPath</key>
    <string>${PWD}/logs/windsurf_integration_output.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>${PWD}</string>
    </dict>
</dict>
</plist>